# mglru0.service.in
[Unit]
Description=Disable mg-LRU
Documentation=https://github.com/hakavlad/mg-lru-helper

[Service]
Type=oneshot
ExecStart=-:TARGET_BIN:/mglru0

[Install]
WantedBy=multi-user.target
